# rpm/noxext.cron.daily
#! /bin/sh

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

tmpreaper 90d --mtime --all /var/log/snacdb
